# node_exporter.service
[Unit]
Description=Node Exporter
After=network-online.target
[Service]
User=promuser
Group=promuser
Type=simple
ExecStart=/usr/local/bin/node_exporter
[Install]
WantedBy=multi-user.target